# packages/thinksuit-tty/etc/thinksuit-tty.service.plist
<!-- Edit below as needed and save to ~/Library/LaunchAgents/thinksuit-tty.service.plist -->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>thinksuit-tty.service</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/mac/Library/Logs/thinksuit-tty.service.stderr.log</string>
    <key>StandardOutPath</key>
    <string>/Users/mac/Library/Logs/thinksuit-tty.service.stdout.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>FORCE_COLOR</key>
      <string>true</string>
      <key>PATH</key>
      <string><![CDATA[/Users/mac/.local/share/mise/installs/node/22/bin:/opt/homebrew/bin:/usr/local/bin:/usr/local/sbin:/usr/bin:/bin:/usr/sbin:/sbin]]></string>
      <key>TTW_PORT</key>
      <string>60662</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/mac/repos/github.com/machellerogden/thinksuit/packages/thinksuit-tty</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/mac/.local/share/mise/installs/node/22/bin/node</string>
      <string>bin/service.mjs</string>
    </array>
  </dict>
</plist>
